# ca-certificates_4.0.1_darwin-amd64/DEBIAN/postinst
#!/bin/bash
launchctl load /Library/LaunchAgents/us.diatr.openssl-osx-ca.plist
/usr/local/bin/openssl-osx-ca -path /usr/local/bin/osx-ca-certs /usr/local/bin/apt